# 99-razer-pwm-fan.rules
# Force-bind Razer PWM PC Fan Controller (1532:0f3c) to our HID driver.
# Note: writing driver_override does NOT unbind automatically; we unbind first, then bind.
# See kernel docs about driver binding/driver_override.

ACTION=="add", SUBSYSTEM=="hid", KERNEL=="0003:1532:0F3C.*", \
  RUN+="/bin/sh -c 'dev=/sys%p; \
                    if [ -L \"$dev/driver\" ]; then echo -n %k > \"$dev/driver/unbind\"; fi; \
                    echo razer_pwm_pc_fan_controller > \"$dev/driver_override\"; \
                    echo -n %k > /sys/bus/hid/drivers/razer_pwm_pc_fan_controller/bind'"
